# feed-irc-bot.service
[Unit]
Description=Feed IRC Bot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /path/to/feed-irc-bot/index.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
